# root/etc/services.d/mmonit/run
#!/usr/bin/with-contenv bash

# Reduce the size of the FD table
ulimit -n ${MMONIT_LIMIT_FD}

# Drop privileges and run M/Monit
s6-setuidgid abc /config/mmonit-${MMONIT_VERSION}/bin/mmonit -i

